# Experiments/systemd/exp01/HelloWorld.service.default
[Unit]
Description=Look Ma No Quotation Marks 
After=network.target

[Service]
Type=simple
ExecStart={{ thisdir }}/script.sh
Restart=always
RestartSec=5s


WantedBy=default.target
#WantedBy=multi-user.target
